# var/spool/cron/crontabs/njoppi2
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.GUfamD/crontab installed on Wed Aug  2 01:24:13 2023)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command

SHELL=/bin/bash
MAILTO=njoppi2@gmail.com
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:/usr/local/bin/myScripts:/usr/local/bin/myScripts/smartBulbScripts:/usr/local/bin/myScripts/screenshotScripts:/usr/local/bin/myScripts/makeComputerUnusable:/usr/local/bin/myScripts/webhooks

# the line below sends your current ip to duckdns every 5m
*/5 * * * * duck.sh >/dev/null 2>&1

# minimizes everything and turn off wifi
30-59 21 * * * makeComputerUnusable &>/tmp/meuCronMakeComputer.log
0 22 * * * makeComputerUsable &>/tmp/meuCronMakeComputer.log

# commit screenshots / print every couple of days
# cant execute the code below because my user (njoppi2) has no permission to write to /tmp (actually i think this only happened because i was trying to rewrite something a /tmp file that was created by root
# 52 * * * 7 /usr/local/bin/myScripts/commitDailyPrints &>/tmp/meuCronRootCommitScreenshots.log
0 8-23 * * 7 commitDailyPrints &>/tmp/meuCronCommitScreenshots.log

#who is crontab here?                         ANSWER: njoppi2
#* * * * * whoami &> /tmp/crontabWhoisThis

# takes a screenshot every 20m
*/20 * * * * perl -le 'sleep rand 150' && takeScreenshotAndSave complete &>/tmp/meuCronScreenshot.log

# changes the light in my room
# I want to go to sleep at 1h30, so 2h before i should start dimming my lights
* * * * * callSmartBulb &>>/tmp/meuCronCallSmartBulb.log

#leave a blank line below this
